# LBG_IA_MMO/infra/systemd/lbg-orchestrator.service
[Unit]
Description=LBG_IA_MMO Orchestrator
After=network.target lbg-agent-dialogue.service

[Service]
Type=simple
User=lbg
Group=lbg
WorkingDirectory=/opt/LBG_IA_MMO/orchestrator
# Métriques GET /metrics (Prometheus texte, opt-in) : défaut désactivé.
# Surcharge : /etc/lbg-ia-mmo.env — LBG_METRICS_ENABLED=1 et optionnellement LBG_METRICS_TOKEN.
Environment=LBG_METRICS_ENABLED=0
EnvironmentFile=-/etc/lbg-ia-mmo.env
Environment=PYTHONUNBUFFERED=1
# Appel POST http://127.0.0.1:8020/invoke pour agent.dialogue (service lbg-agent-dialogue)
Environment=LBG_AGENT_DIALOGUE_URL=http://127.0.0.1:8020
ExecStart=/opt/LBG_IA_MMO/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8010
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target

